# Zillow_House_Info_Clawer/skadihouse.cron
52 17 */3 * * (/bin/sh /Spring-Cloud-SkadiHouse/Zillow_House_Info_Clawer/scrapy_start.sh)  
