# Zillow_House_Info_Clawer/skadihouse.cron
PATH=/usr/bin
* 12 * * * cd /home/vagrant/workspace/Spring-Cloud-SkadiHouse/Zillow_House_Info_Clawer/ && scrapy crawl skadihouse
